# deploy/enterprise/fs/etc/cron.d/renew-certificate
* * * * 0 root exec /opt/steedos/renew-certificate.sh